# host/misc/udev/88-nuand.rules
# nuand bladeRF
ATTR{idVendor}=="1d50", ATTR{idProduct}=="6066", MODE="660", GROUP="adm", RUN+="/usr/share/bladerf/scripts/init_board.sh"
